# tools/launchd/rust-mux.sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.libraxis.mcp-mux.general-memory</string>

  <!-- Adjust paths below to your environment -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USER/.cargo/bin/rust-mux</string>
    <string>--config</string>
    <string>/Users/YOUR_USER/.codex/mcp.json</string>
    <string>--service</string>
    <string>general-memory</string>
    <string>--tray</string>
  </array>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/rust-mux.general-memory.out</string>
  <key>StandardErrorPath</key><string>/tmp/rust-mux.general-memory.err</string>
</dict>
</plist>
